# dnssec-triggerd-keygen.service.in
[Unit]
Description=Generate keys and certificates for dnssec-trigger
ConditionPathExists=!@sysconfdir@/dnssec_trigger_control.key

[Service]
Type=oneshot
ExecStart=/usr/sbin/dnssec-trigger-control-setup -d @sysconfdir@
RemainAfterExit=yes
